# files/etc/uci-defaults/cloudflared
#!/bin/bash
# Enable Cloudflare Tunnel default
uci set cloudflared.config.enabled=1
uci set cloudflared.config.token=''
uci set cloudflared.config.region='us'
uci commit
