# files/etc/uci-defaults/upnp
#!/bin/bash
uci set upnpd.config.enabled=1
uci commit
